# docker/php/crontab
* * * * * cd /var/www/src && /usr/local/bin/php artisan schedule:run >> /var/log/cron.log 2>&1
